# packaging/systemd/olcrtc-tor-bridge-pool.service
[Unit]
Description=Probe Tor bridge pool (webtunnel/obfs4) and refresh bridges.conf
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@OLC_SCRIPTS@/tor-bridge-pool.sh --url-only --jobs 8 --target 12 --types webtunnel,obfs4
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
