# packaging/systemd/olcrtc-tor-bridge-monitor.timer
[Unit]
Description=Tor bridge health monitor (every 20 min)

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
RandomizedDelaySec=90s
Persistent=true

[Install]
WantedBy=timers.target
